# files/etc/init.d/wanmonitor
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/bin/wanmonitor -d &
}

stop() {
	kill $(pidof wanmonitor)
}
